# deploy/systemd/tiyf-engine-demo.service
[Unit]
Description=TiYf Engine (demo)
Wants=network-online.target
After=network-online.target

[Service]
# Run as root (only user on VPS)
WorkingDirectory=/opt/tiyf/current

# systemd reads the environment file as root before starting
EnvironmentFile=/etc/tiyf/engine.env

# Safe defaults (only used if file missing)
Environment=ADAPTER_MODE=oanda-demo
Environment=ENGINE_CONFIG_PATH=/opt/tiyf/current/sample-config.demo-oanda.json

# Launch directly; no manual file reads or bash sourcing needed
ExecStart=/usr/bin/dotnet /opt/tiyf/current/TiYf.Engine.Host.dll --config ${ENGINE_CONFIG_PATH}

Restart=always
RestartSec=3

# Optional hardening
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/opt/tiyf /etc/tiyf
PrivateTmp=yes

[Install]
WantedBy=multi-user.target